# scripts/webkassa-key-update.service
[Unit]
Description=Update Webkassa API Key
After=network.target docker.service
Requires=docker.service

[Service]
Type=oneshot
User=root
WorkingDirectory=/opt/altegio-webkassa-integration
ExecStart=/usr/local/bin/update-webkassa-key.sh
StandardOutput=journal
StandardError=journal

# Restart policy
Restart=no
TimeoutStopSec=30

# Environment
Environment=DOCKER_COMPOSE_FILE=/opt/altegio-webkassa-integration/docker-compose.yml

[Install]
WantedBy=multi-user.target
