# scripts/webkassa-key-update.timer
[Unit]
Description=Run Webkassa API Key Update Daily at 5 AM
Requires=webkassa-key-update.service

[Timer]
# Запуск каждый день в 5:00 утра
OnCalendar=*-*-* 05:00:00
# Если система была выключена, запустить при следующем включении
Persistent=true
# Добавить случайную задержку до 30 минут для избежания нагрузки
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target
